# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :99 -screen 0 1280x1024x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
environment=DISPLAY=:99

[program:x11vnc]
command=x11vnc -display :99 -forever -nopw -shared -passwd secret
autostart=true
autorestart=true

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true